# Profile/airootfs/home/arcade/.config/autostart/calamares.desktop
[Desktop Entry]
Name=Install Arcade Linux
Exec=sudo calamares
Type=Application
Terminal=false
